# root/etc/service/php-fpm/run
#!/bin/bash

mkdir -p /run/php && chown -R app:app /run/php
exec /usr/sbin/php-fpm${PHP_VERSION} -F --pid=/tmp/php-fpm.pid

